# scripts/systemd/ntfy-heartbeat-check.timer
[Unit]
Description=Run bot heartbeat check every 3 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
