# scripts/systemd/ntfy-disk-check.timer
[Unit]
Description=Check disk usage every 6 hours

[Timer]
OnBootSec=15min
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
